# .config/systemd/user/bat.service
[Unit]
Description=Run bat-stat.sh

[Service]
Type=oneshot
ExecStart=/usr/bin/sh -c '$HOME/bin/bat-stat.sh'
