# .config/systemd/user/bat.service
[Unit]
Description=Run bat-stat.sh

[Service]
Type=oneshot
ExecStart=/usr/bin/sh -c '/home/lam/bin/bat-stat.sh'
